# agent/scripts/iotistic-agent.service
[Unit]
Description=Iotistic Agent - IoT Device Management Service
Documentation=https://github.com/Iotistica/iotistic
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=iotistic
Group=iotistic
WorkingDirectory=/opt/iotistic/agent

# Environment
EnvironmentFile=/etc/iotistic/agent.env

# Start command
ExecStart=/usr/bin/node dist/index.js

# Restart behavior
Restart=always
RestartSec=10
StartLimitInterval=400
StartLimitBurst=3

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
SendSIGKILL=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/iotistic/agent /var/log/iotistic

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iotistic-agent

[Install]
WantedBy=multi-user.target
